# ds4drv/debian/postinst
#!/bin/sh
#------------------------------------------------------------------------------------------------------
# Help notes:
# See: https://www.debian.org/doc/manuals/debian-faq/ch-pkg_basics.en.html
#------------------------------------------------------------------------------------------------------
# This script typically completes any required configuration of the package foo once foo has been 
# unpacked from its Debian archive (".deb") file. Often, 'postinst' scripts ask the user for input, 
# and/or warn the user that if he accepts default values, he should remember to go back and re-configure
# that package as the situation warrants. Many 'postinst' scripts then execute any commands necessary to
# start or restart a service once a new package has been installed or upgraded.
#
# see: dh_installdeb(1)

set -e


case "$1" in
    configure)

	# enable/start the ds4drv service
	# Start systemd services
#	echo "Staring ds4drv service..."
#	sudo systemctl enable ds4drv.service
#	sudo systemctl start ds4drv.service
    :
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

